# roles/cca.core.linux/templates/etc/systemd/system/cca_cgroup_version_daemon.service.j2
#This daemon is used to make sure cgroupv1 is being used and is available
#even after a manual kernel update.
#systemd-based Linux distributions.
#Managed by CCA as a template file in Ansible

[Unit]
Description=Systemd service file for CCA cgroup version daemon, generated by CCA
After=network.target
StartLimitInterval=200
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/cca_cgroup_version_handler.sh

[Install]
WantedBy=multi-user.target